# notify-panel.service
[Unit]
Description=Telegram Notify Panel
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=/opt/notify-panel
ExecStart=/usr/bin/python3 /opt/notify-panel/app.py
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

